# quiz_app_supervisor.conf
[program:quick_app_backend]
directory=/home/ubuntu/quick_app_backend
command=/bin/bash -E -c ./backend_start.sh
autostart=true
autorestart=true
stopsignal=INT
stopasgroup=true
killasgroup=true